# asset/etc/crontabs/root
* * * * * php /data/wwwroot/artisan schedule:run >> /dev/null 2>&1